# roles/restic/templates/restic.timer.j2
[Unit]
Description=Run {{ restic_job_name }}

[Timer]
OnCalendar={{ restic_policy.frequency }}
Persistent=True
Unit={{ restic_systemd_unit_naming_scheme }}.service

[Install]
WantedBy=timers.target
